# doc/examples/systemd/mlvpn@.service.in
[Unit]
Description=MLVPN connection to %i
PartOf=mlvpn.service
ReloadPropagatedFrom=mlvpn.service

[Service]
Type=notify
NotifyAccess=main
ExecStart=@sbindir@/mlvpn --config @sysconfdir@/mlvpn/%i.conf --name %i --user mlvpn
ExecReload=/bin/kill -HUP $MAINPID
WorkingDirectory=@sysconfdir@/mlvpn
Restart=on-failure

[Install]
WantedBy=multi-user.target


